# telegram-gateway.service
[Unit]
Description=Telegram Gateway for Gemini CLI
After=network.target

[Service]
Type=simple
User=gemini
# Измените путь на актуальный для вашей системы
WorkingDirectory=/opt/telegram-gateway
EnvironmentFile=/opt/telegram-gateway/.env
# Путь к python в виртуальном окружении
ExecStart=/opt/telegram-gateway/.venv/bin/python -m gateway.main
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
