# telegram-gateway.service
# Template for manual installation or for install.sh placeholder substitution.
# Replace every __PLACEHOLDER__ value before copying this file to /etc/systemd/system/.

[Unit]
Description=Telegram Gateway for Gemini CLI
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=__ENV_FILE__
Environment=HOME=__HOME_DIR__
Environment=PYTHONUNBUFFERED=1
Environment=PATH=__SERVICE_PATH__
ExecStart=__PYTHON_BIN__ -m gateway.main
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
